# systemd/palworld-panel.service
[Unit]
Description=Palworld Server Management Panel
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/palworld-panel/backend
Environment="PATH=/opt/palworld-panel/backend/venv/bin"
ExecStart=/opt/palworld-panel/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target